# backend/deploy/zup.service
[Unit]
Description=Zup Backend (Flask + Socket.IO)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/zup/backend
EnvironmentFile=-/var/www/zup/backend/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/var/www/zup/backend/venv/bin/gunicorn \
    --worker-class gthread \
    --workers 1 \
    --threads 4 \
    --bind 127.0.0.1:5003 \
    --timeout 120 \
    wsgi:app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
